# cpu_fan.service
[Unit]
Description=Fan Control

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/cpu_fan

[Install]
WantedBy=multi-user.target
